# bin/noderunner.service
[Unit]
Description=Daemon for running jobs from Mongo queue

[Service]
ExecStart=/usr/bin/node noderunner.js
WorkingDirectory=/home/noderunner
Restart=always
RestartSec=500ms
StartLimitInterval=0
Environment=ENV=production
User=nginx
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=noderunner

[Install]
WantedBy=multi-user.target
